# configs/supervisord.conf
[supervisord]
nodaemon=true

[program:seafile]
command=/bin/seafile-wrapper.sh 
[program:seahub]
command=/bin/seahub-wrapper.sh 
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
